# configs/udev/99-hdmi-audio-fix.rules
ACTION=="change", SUBSYSTEM=="drm", ENV{HOTPLUG}=="1", RUN+="/home/max/.local/bin/hdmi-audio-hotplug.sh"
